# config/scraper.service
[Unit]
Description=Web scraper and sender

[Service]
ExecStart=python3 /opt/scraper/main.py